# profile.d/custom.sh
#!/bin/bash
export PATH=$PATH:/home/andrew/.bin/
